# bin/BrokerBot.service
[Unit]
Description=Lab-93 Broker Bot Daemon
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=hunter
WorkingDirectory=/server/finance/Broker-Bot
VIRTUAL_ENV=/server/finance/Broker-Bot/
ExecStart=python3 BrokerBot.py

[Install]
WantedBy=multi-user.target
